# deploy/polybot.service
# systemd unit template for the Polymarket BTC 15-min bot.
#
# Copy to /etc/systemd/system/polybot.service, edit the paths if your install
# isn't at /opt/polybot, then:
#
#   sudo systemctl daemon-reload
#   sudo systemctl enable polybot
#   sudo systemctl start polybot
#   journalctl -u polybot -f          # watch startup
#
# IMPORTANT — Restart policy: this unit uses ``Restart=no`` so the bot does NOT
# auto-restart after a fail-stop. The ledger design fail-stops on durable-
# write failure; restarting silently would weaken the manual-reconciliation
# guarantee. After a failure, the operator must inspect the journal, reconcile
# any unresolved ledger state, and then ``systemctl start polybot`` explicitly.
#
# Live mode decrypts credentials/encrypted_credentials.json. The service asks
# systemd for the vault password as root, then runs bot.py as the polybot user.

[Unit]
Description=Polymarket BTC 15-Min Trading Bot
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polybot/Polymarket-BTC-15-Minute-Trading-Bot

Environment=LIVE_TRADE_LEDGER_PATH=/opt/polybot/ledger/live_trades.json
Environment=NAUTILUS_LOG_DIR=/opt/polybot/logs/nautilus
Environment=DECISION_LOG_PATH=/opt/polybot/ledger/decisions.jsonl

ExecStart=/bin/sh -c 'systemd-ask-password "Polymarket credentials vault password:" | /usr/sbin/runuser --preserve-environment -u polybot -- /opt/polybot/venv/bin/python bot.py --live --confirm-live'

Restart=no

# Resource limits — modest; the bot is I/O bound.
MemoryMax=1G
TasksMax=200

# Logging to journald (paired with logrotate)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polybot

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/polybot/ledger /opt/polybot/logs /run/systemd/ask-password
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
